# conf/project/circus/circus.conf
{% import '_vars.sls' as vars with context %}
{% set env_dir = vars.path_from_root('env') %}

start on filesystem and net-device-up IFACE=lo

stop on shutdown

respawn
exec {{ env_dir }}/bin/circusd --log-output {{ vars.log_dir }}/circus.log \
                            --pidfile /var/run/circusd.pid \
                            /etc/circus.ini
